# dpkg/etc/systemd/system/stund.service
[Unit]
Description=STUN for BerryShare
After=network.target

[Service]
ExecStart=/usr/bin/stund
Restart=always

[Install]
WantedBy=multi-user.target
